# packages/active-listener/systemd/active-listener.service.in
[Unit]
Description=Eavesdrop active-listener desktop dictation service
PartOf=graphical-session.target
After=graphical-session.target
StartLimitIntervalSec=1800
StartLimitBurst=60

[Service]
Type=simple
WorkingDirectory=@REPO_ROOT@
ExecStart=/bin/bash -c 'set -e; ENV=$(timeout 2 cat @REPO_ROOT@/.env); set -a; eval "$ENV"; exec @UV_BIN@ run --project @REPO_ROOT@/packages/active-listener active-listener'
Restart=on-failure
RestartSec=30

[Install]
WantedBy=graphical-session.target
